# ansible/dist/etc/supervisor/conf.d/kibana.conf
[program:kibana]
command=/opt/kibana/bin/kibana
autostart=true
autorestart=true
stderr_logfile=/var/log/kibana.err.log
stdout_logfile=/var/log/kibana.out.log